# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
; Increase request/header buffers to handle large Cookie headers
; and prevent "invalid request block size: ... (max 4096)" that surfaces as 502 at proxy.
command=/bin/sh -c "python manage.py migrate_schemas --noinput && python manage.py collectstatic --noinput && uwsgi --http :8002 --module config.wsgi:application --master --processes 4 --threads 2 --buffer-size 65535 --http-buffer-size 65535"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
