# supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:uwsgi]
; Increase buffer sizes to handle large Authorization/Cookie headers
; and prevent "invalid request block size" (502) at reverse proxy.
command=/bin/sh -c "python manage.py migrate_schemas --noinput && python manage.py seed_tenant && python manage.py collectstatic --noinput && uwsgi --http :8001 --module config.wsgi:application --master --processes 4 --threads 2 --buffer-size 65535 --http-buffer-size 65535 --harakiri 60 --max-requests 1000"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
