# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=/bin/sh -c "python manage.py migrate_schemas --noinput && python manage.py collectstatic --noinput && gunicorn config.wsgi:application --bind 0.0.0.0:8002 --workers 3 --timeout 120 --access-logfile - --error-logfile -"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
